# test/file.conf
# [program:mywebapp]
# command=/home/hasabir/miniconda3/bin/python /home/hasabir/Desktop/projects/taskmaster/test/app.py
# directory=/home/hasabir/Desktop/projects/taskmaster/test
# autostart=true
# autorestart=true
# startretries=3
# user=hasabir  # Use your actual username, not www-data
# environment=HOME="/home/hasabir",USER="hasabir",PYTHONPATH="/home/hasabir/Desktop/projects/taskmaster/test"
# stdout_logfile=/home/hasabir/Desktop/projects/taskmaster/test/mywebapp.out.log
# stderr_logfile=/home/hasabir/Desktop/projects/taskmaster/test/mywebapp.err.log
# stdout_logfile_maxbytes=10MB
# stdout_logfile_backups=5



; /etc/supervisor/conf.d/myapps.conf
; Multiple programs configuration

; === Web Application ===





; === Background Worker ===
[program:worker]
command=/home/hasabir/miniconda3/bin/python /home/hasabir/Desktop/projects/taskmaster/test/worker.py
directory=/home/hasabir/Desktop/projects/taskmaster/testfere
autostart=false
# autorestart=true
startretries=3
user=hasabir
environment=HOME="/home/hasabir",USER="hasabir"
# stdout_logfile=/home/hasabir/Desktop/projects/taskmaster/test/worker.out.log
# stderr_logfile=/home/hasabir/Desktop/projects/taskmaster/test/worker.err.log
exitcodes=0





; === API Server (different port) ===
# [program:api]
# command=/home/hasabir/miniconda3/bin/python /home/hasabir/Desktop/projects/taskmaster/test/api_server.py
# directory=/home/hasabir/Desktop/projects/taskmaster/test
# autostart=true
# autorestart=true
# startretries=3
# numprocs=3
# process_name=%(program_name)s_%(process_num)02d
# user=hasabir
# environment=HOME="/home/hasabir",USER="hasabir",PORT=5001
# stdout_logfile=/home/hasabir/Desktop/projects/taskmaster/test/api.out.log
# stderr_logfile=/home/hasabir/Desktop/projects/taskmaster/test/api.err.log

# ; === Health Monitor ===
# [program:monitor]
# command=/home/hasabir/miniconda3/bin/python /home/hasabir/Desktop/projects/taskmaster/test/monitor.py
# directory=/home/hasabir/Desktop/projects/taskmaster/test
# autostart=true
# autorestart=true
# startretries=3
# user=hasabir
# environment=HOME="/home/hasabir",USER="hasabir"
# stdout_logfile=/home/hasabir/Desktop/projects/taskmaster/test/monitor.out.log
# stderr_logfile=/home/hasabir/Desktop/projects/taskmaster/test/monitor.err.log

# [program:worker_multi]
# command=/home/hasabir/miniconda3/bin/python /home/hasabir/Desktop/projects/taskmaster/test/worker_multi.py
# process_name=%(program_name)s_%(process_num)02d
# numprocs=3
# numprocs_start=0
# directory=/home/hasabir/Desktop/projects/taskmaster/test
# autostart=true
# autorestart=true
# user=hasabir

# ; === YOUR SPECIFIED CONFIGURATIONS ===
# exitcodes=0                    ; Consider exit code 0 as success
# startretries=3                 ; Retry 3 times before giving up
# startsecs=5                    ; Consider start successful after 5 seconds
# stopsignal=USR1                ; Use USR1 signal for graceful shutdown
# stopwaitsecs=10                ; Wait 10 seconds for graceful shutdown

# ; Logging
# stdout_logfile=/home/hasabir/Desktop/projects/taskmaster/test/worker_multi_%(process_num)02d.out.log
# stderr_logfile=/home/hasabir/Desktop/projects/taskmaster/test/worker_multi_%(process_num)02d.err.log
# stdout_logfile_maxbytes=10MB
# stdout_logfile_backups=5

# ; === Web Application ===
# [program:webapp]
# command=/home/hasabir/miniconda3/bin/python /home/hasabir/Desktop/projects/taskmaster/test/app.py
# directory=/home/hasabir/Desktop/projects/taskmaster/test
# autostart=true
# autorestart=true
# user=hasabir
# stdout_logfile=/home/hasabir/Desktop/projects/taskmaster/test/webapp.out.log
# stderr_logfile=/home/hasabir/Desktop/projects/taskmaster/test/webapp.err.log